# packages/mail-mta/opensmtpd/files/systemd/smtpd.service
[Unit]
Description=SMTP server

[Service]
ExecStart=/usr/sbin/smtpd -d

[Install]
WantedBy=multi-user.target
