# packages/www-servers/nginx/files/systemd/nginx.service
[Unit]
Description=nginx HTTP server

[Service]
Type=forking
ExecStart=/usr/sbin/nginx
ExecReload=/usr/sbin/nginx -s reload
ExecStop=/usr/sbin/nginx -s stop

[Install]
WantedBy=network.target
