# packages/app-admin/puppet/files/systemd/puppet-master.service
[Unit]
Description=puppet master daemon

[Service]
ExecStart=/usr/bin/puppet master --no-daemon
Restart=always

[Install]
WantedBy=multi-user.target
