# packages/net-dns/unbound/files/systemd/unbound.service
[Unit]
Description=unbound caching DNS resolver

[Service]
ExecStart=/usr/sbin/unbound -d
Restart=always

[Install]
WantedBy=network.target
